# dist/spuddns.service
[Unit]
Description=spuddns DNS Resolver
After=network.target

[Service]
ExecStart=/opt/spuddns /etc/spuddns.json
Restart=always
ExecStop=/bin/kill -s QUIT $MAINPID
TimeoutStopSec=6

[Install]
Alias=spuddns.service
WantedBy=multi-user.target
